# files/etc/init.d/performance-mode
#!/bin/sh /etc/rc.common

START=95

# 返回可用于网络队列分流的 CPU 核心数。
get_cpu_count() {
	local count

	count="$(nproc 2>/dev/null || true)"
	case "$count" in
		''|*[!0-9]*)
			count="$(grep -c '^processor' /proc/cpuinfo 2>/dev/null || echo 1)"
			;;
	esac

	case "$count" in
		''|*[!0-9]*|0)
			count=1
			;;
	esac

	echo "$count"
}

# 按 CPU 核心数生成十六进制掩码，供 RPS/XPS 队列绑定使用。
build_cpu_mask() {
	local count="$1"
	local chunk_hex
	local mask=""

	while [ "$count" -gt 0 ]; do
		if [ "$count" -ge 32 ]; then
			chunk_hex="ffffffff"
			count=$((count - 32))
		else
			chunk_hex="$(printf '%x' $(( (1 << count) - 1 )))"
			count=0
		fi

		if [ -n "$mask" ]; then
			mask="${chunk_hex},${mask}"
		else
			mask="$chunk_hex"
		fi
	done

	echo "$mask"
}

# 把 CPU governor 和 EPP 拉到性能档，减少代理转发高峰时的突发降频。
apply_cpu_governor() {
	for governor_file in \
		/sys/devices/system/cpu/cpufreq/policy*/scaling_governor \
		/sys/devices/system/cpu/cpu*/cpufreq/scaling_governor; do
		[ -f "$governor_file" ] || continue

		available_file="${governor_file%/*}/scaling_available_governors"
		if [ -f "$available_file" ] && grep -qw performance "$available_file"; then
			echo performance >"$governor_file" 2>/dev/null || true
		fi
	done

	for epp_file in \
		/sys/devices/system/cpu/cpufreq/policy*/energy_performance_preference \
		/sys/devices/system/cpu/cpu*/cpufreq/energy_performance_preference; do
		[ -f "$epp_file" ] || continue

		available_file="${epp_file%/*}/energy_performance_available_preferences"
		if [ -f "$available_file" ] && grep -qw performance "$available_file"; then
			echo performance >"$epp_file" 2>/dev/null || true
		fi
	done
}

# 提前加载 TCP 拥塞控制模块，让 Turbo ACC 或 sysctl 能看到 BBR/Hybla/Scalable。
load_tcp_congestion_modules() {
	local module

	for module in tcp_bbr tcp_hybla tcp_scalable; do
		modprobe "$module" >/dev/null 2>&1 || true
	done
}

# 调整网络栈参数，照顾 SSR-Plus 本机代理连接和 DNS 服务的突发请求。
apply_sysctl_tweaks() {
	load_tcp_congestion_modules

	sysctl -w net.core.default_qdisc=fq >/dev/null 2>&1 || true
	sysctl -w net.core.somaxconn=4096 >/dev/null 2>&1 || true
	sysctl -w net.ipv4.tcp_congestion_control=bbr >/dev/null 2>&1 || true
	sysctl -w net.ipv4.tcp_fastopen=3 >/dev/null 2>&1 || true
	sysctl -w net.ipv4.tcp_mtu_probing=1 >/dev/null 2>&1 || true
	sysctl -w net.ipv4.tcp_slow_start_after_idle=0 >/dev/null 2>&1 || true
	sysctl -w net.ipv4.tcp_max_syn_backlog=4096 >/dev/null 2>&1 || true
	sysctl -w net.ipv4.tcp_fin_timeout=15 >/dev/null 2>&1 || true
	sysctl -w net.core.netdev_max_backlog=8192 >/dev/null 2>&1 || true
	sysctl -w net.core.netdev_budget=600 >/dev/null 2>&1 || true
	sysctl -w net.core.netdev_budget_usecs=8000 >/dev/null 2>&1 || true
}

# 为每个网卡队列设置 RPS/XPS 和 flow entries，降低单核心软中断压力。
apply_queue_tuning() {
	local cpu_count cpu_mask iface rx_queue flow_cnt rx_queue_count

	cpu_count="$(get_cpu_count)"
	[ "$cpu_count" -gt 1 ] || return 0

	cpu_mask="$(build_cpu_mask "$cpu_count")"
	sysctl -w net.core.rps_sock_flow_entries=32768 >/dev/null 2>&1 || true

	for iface in /sys/class/net/*; do
		[ -d "$iface" ] || continue

		case "${iface##*/}" in
			lo)
				continue
				;;
		esac

		rx_queue_count=0
		for rx_queue in "$iface"/queues/rx-*; do
			[ -f "$rx_queue/rps_cpus" ] || continue
			echo "$cpu_mask" >"$rx_queue/rps_cpus" 2>/dev/null || true
			rx_queue_count=$((rx_queue_count + 1))
		done

		if [ "$rx_queue_count" -gt 0 ]; then
			flow_cnt=$((32768 / rx_queue_count))
			[ "$flow_cnt" -ge 256 ] || flow_cnt=256

			for rx_queue in "$iface"/queues/rx-*; do
				[ -f "$rx_queue/rps_flow_cnt" ] || continue
				echo "$flow_cnt" >"$rx_queue/rps_flow_cnt" 2>/dev/null || true
			done
		fi

		for tx_queue in "$iface"/queues/tx-*; do
			[ -f "$tx_queue/xps_cpus" ] || continue
			echo "$cpu_mask" >"$tx_queue/xps_cpus" 2>/dev/null || true
		done
	done
}

# 启动时按顺序应用 CPU、网络栈、队列分流三类优化。
start() {
	apply_cpu_governor
	apply_sysctl_tweaks
	apply_queue_tuning
}

# 重启服务时直接复用启动逻辑，保证参数重新写回。
restart() {
	start
}
