# files/etc/init.d/ssr-optimization
#!/bin/sh /etc/rc.common
# SSR-Plus Performance Optimization Script
# Priority: 95 (run after network but before apps)

START=95
STOP=15
USE_PROCD=1

start() {
    echo "Starting SSR-Plus optimizations..."
    
    # CPU频率设置（性能模式）
    echo performance > /sys/devices/system/cpu/cpu0/cpufreq/scaling_governor 2>/dev/null
    
    # 网络队列优化
    for nic in $(ls /sys/class/net/ | grep -E 'eth|wlan'); do
        echo 1024 > /sys/class/net/$nic/tx_queue_len 2>/dev/null
        ethtool -K $nic gso on gro on tso on 2>/dev/null
    done
    
    # 中断亲和性（绑定到所有CPU）
    echo ff > /proc/irq/default_smp_affinity 2>/dev/null
    
    # 增加文件描述符限制
    ulimit -n 65535
    ulimit -u 32768
    
    # 启用TCP BBR（如果内核支持）
    modprobe tcp_bbr 2>/dev/null
    echo bbr > /proc/sys/net/ipv4/tcp_congestion_control 2>/dev/null
    
    # 启用Flow Offload
    echo 1 > /proc/sys/net/netfilter/nf_flowtable_tcp_timeout 2>/dev/null
    echo 1 > /proc/sys/net/netfilter/nf_flowtable_udp_timeout 2>/dev/null
    
    # DNS缓存优化
    echo 8192 > /proc/sys/net/core/netdev_max_backlog 2>/dev/null
    
    # 内存优化
    sync && echo 3 > /proc/sys/vm/drop_caches 2>/dev/null
    
    echo "SSR-Plus optimizations applied."
}

stop() {
    echo "SSR-Plus optimizations remain active."
}

restart() {
    start
}
